# ansible/roles/pve_host/files/nvme-metrics.service
[Unit]
Description=Collect NVMe SMART metrics for node-exporter

[Service]
Type=oneshot
ExecStart=/usr/local/bin/nvme-metrics
